# example-run-scripts/gufi-update.timer
[Unit]
Description=Gufi reindex list timer

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
